# systemd/ibeach-notify@.service
[Unit]
Description=%i failure email notification

[Service]
Type=simple
ExecStart=/home/pi/git/ipcam/systemd/email_notifier.sh
